# etc/init/me135.conf
exec /usr/local/bin/node /home/stvn/me135/beaglebone/main.js
start on runlevel 2
respawn
